# saori-revive/saori-bot.service
[Unit]
Description=Saori Discord Music Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=warwick
Group=warwick
WorkingDirectory=/home/warwick/Saori_revive
Environment=PYTHONPATH=/home/warwick/Saori_revive
ExecStart=/home/warwick/Saori_revive/.venv/bin/python main.py
Restart=on-failure
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/warwick/Saori_revive

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=saori-bot

[Install]
WantedBy=multi-user.target